# dot_crontab.cron
0 12 * * 2 ~/.local/scripts/git-afforester worktree prune
